# dnd-bot.service
[Unit]
Description=Bot de D&D 5e para WhatsApp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Ajuste o usuário e o caminho para onde você copiou o projeto na VM
User=ubuntu
WorkingDirectory=/home/ubuntu/dnd5e-whatsapp-bot
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/ubuntu/dnd5e-whatsapp-bot/.env

[Install]
WantedBy=multi-user.target
